# scripts/collectSMART_v2.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games

@hourly /home/zfs/SMART_cron/checkSMARThourly_v2.sh



